# etc/systemd/system/restic-prune.service
[Unit]
Description=Restic Prune service

[Service]
Type=oneshot
ExecStart=restic-offsite prune
